# mac/LaunchAgents/com.ostepan.display-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ostepan.display-watcher</string>

    <!-- display-watcher polls aerospace list-monitors every ~3s and runs
         set-main-display.sh auto whenever the layout changes. Polling is
         used instead of CGDisplayRegisterReconfigurationCallback because
         CG callbacks silently miss events from a launchd-context process. -->
    <!-- $HOME-relative so it's portable across machines (launchd sets HOME). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec "$HOME/.config/aerospace/display-watcher.sh"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if it ever exits — RunLoop.current.run() never returns in
         normal operation, but if the binary segfaults or is killed we want
         it back. -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/aerospace-display-watcher.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/aerospace-display-watcher.err</string>
</dict>
</plist>
